# deployment/mercedes-finder.service
[Unit]
Description=Mercedes W123/W124 Diesel Finder
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/mercedes-finder
Environment="PATH=/opt/mercedes-finder/venv/bin"
ExecStart=/opt/mercedes-finder/venv/bin/python main.py
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/mercedes-finder

[Install]
WantedBy=multi-user.target
